# deploy/systemd/classroom-backend.service
[Unit]
Description=Classroom Control Backend
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/classroom/backend
ExecStart=/opt/classroom/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8082
Restart=always
RestartSec=5
EnvironmentFile=/etc/classroom-control/secrets

[Install]
WantedBy=multi-user.target
